# deploy/mealhouse-backup.timer
[Unit]
Description=Nightly MealHouse backup

[Timer]
OnCalendar=*-*-* 02:17:00
Persistent=true
RandomizedDelaySec=900

[Install]
WantedBy=timers.target

